# deploy/widget-ui.service
[Unit]
Description=widget-ui (FireWeave demo)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/fireweave/widget-ui
EnvironmentFile=-/opt/fireweave/widget-ui/.env
Environment=PORT=3102
Environment=HOST=0.0.0.0
ExecStart=/root/.bun/bin/bun ./build/index.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
